# systemd/lowbat.service
[Unit]
Description=Battery discharge watcher

[Service]
Type=oneshot
ExecStart=/home/samo/.local/bin/scripts/lowbat.bash
