# systemd/lowbat.timer
[Unit]
Description=Battery Check Timer

[Timer]
OnActiveSec=1s
OnUnitActiveSec=1m

[Install]
WantedBy=default.target
