# systemd/repomount.service
[Unit]
Description=Portage squashed repository mounting service
After=local-fs.target systemd-tmpfiles-setup.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/squashsync --init

[Install]
WantedBy=default.target
