# systemd/startx.target
[Unit]
Description=X is running (started by xinitrc)
BindsTo=graphical-session.target
